# dispatcher/install/voice-dispatcher.service
[Unit]
Description=voice-dispatcher — ambient voice control for claude-code-hermit
After=network.target sound.target

[Service]
Type=simple
# Set by install-linux.sh
ExecStart=VENV_PYTHON_PLACEHOLDER -m voice_dispatcher run
Restart=on-failure
RestartSec=5
# Mic access — user session service; no special audio group needed with PipeWire/PulseAudio
Environment=DISPLAY=:0

[Install]
WantedBy=default.target
