# dispatcher/install/voice-dispatcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.gtapps.voice-dispatcher</string>

  <key>ProgramArguments</key>
  <array>
    <!-- Absolute path set by install-macos.sh -->
    <string>VENV_PYTHON_PLACEHOLDER</string>
    <string>-m</string>
    <string>voice_dispatcher</string>
    <string>run</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <!-- Log files — created automatically -->
  <key>StandardOutPath</key>
  <string>REPLACE_LOG_DIR/voice-dispatcher.log</string>
  <key>StandardErrorPath</key>
  <string>REPLACE_LOG_DIR/voice-dispatcher.error.log</string>

  <!-- Required for mic access; the process runs in the user session -->
  <key>SessionCreate</key>
  <true/>
</dict>
</plist>
